# packaging/el/inc/files
%files
%defattr(-,root,root,-)

/usr/share/python/sd-agent/agent.py
/usr/share/python/sd-agent/aggregator.py
/usr/share/python/sd-agent/config.py
/usr/share/python/sd-agent/daemon.py
/usr/share/python/sd-agent/emitter.py
/usr/share/python/sd-agent/jmxfetch.py
/usr/share/python/sd-agent/modules.py
/usr/share/python/sd-agent/sd-cert.pem
/usr/share/python/sd-agent/transaction.py
/usr/share/python/sd-agent/util.py

/usr/share/python/sd-agent/checks
/usr/share/python/sd-agent/utils

%config(noreplace) /etc/sd-agent/supervisor.conf
%config(noreplace) /etc/sd-agent/config.cfg
%config(noreplace) /etc/sd-agent/plugins.cfg

/usr/share/python/sd-agent/bin
/usr/share/python/sd-agent/include
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/*.py
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/distutils
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/encodings
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/lib-dynload
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/no-global-site-packages.txt
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/orig-prefix.txt
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/backports*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/boto*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/*consul*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/dns*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/docker*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/easy-install*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/*etcd*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/google
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/ipaddress.py*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/kazoo
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/kazoo-*-info
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/meld3*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/ntplib*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/pip*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/protobuf*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/psutil*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/pycurl*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/PyYAML*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/requests*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/setuptools*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/simplejson*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/six.py*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/six*-info
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/supervisor*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/tornado*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/uptime*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/urllib3*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/websocket*
/usr/share/python/sd-agent/lib/python%{__sd_python_version}/site-packages/yaml

/etc/init.d/sd-agent
/usr/share/python/sd-agent/lib64
